# profile.d/dbod-core.sh
dbod_core_base=/opt/dbod/perl5
export MANPATH=${dbod_core_base}/man:${MANPATH}
export PERL5LIB=${PERL5LIB}:${dbod_core_base}/lib/perl5
export PATH=${PATH}:${dbod_core_base}/bin

